# peapod.service.in
[Unit]
Description=EAPOL Proxy Daemon
Documentation=man:peapod(8) man:peapod.conf(5)
After=network.target

[Service]
ExecStart=__SBIN__/peapod -d -p /var/run/peapod.pid -l /var/log/peapod.log
Type=forking
PIDFile=/var/run/peapod.pid

[Install]
WantedBy=multi-user.target
